# docker-apache-php/conf/supervisord.conf
[supervisord]
nodaemon = true

[program:sshd]
command = /usr/sbin/sshd -D

[include]
files = supervisor/conf.d/*.conf